# deploy/systemd/lunarleague-auto-deploy.timer
[Unit]
Description=Timer periodic Lunar League deploy from git

[Timer]
# Every 2 minutes (homelab / fast iteration). For production, prefer something
# like OnCalendar=*-*-* 00,06,12,18:00 with RandomizedDelaySec=300.
OnCalendar=*-*-* *:0/2:00
Persistent=false
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
