# deploy/systemd/lunarleague-auto-deploy.timer
[Unit]
Description=Timer periodic Lunar League deploy from git

[Timer]
# Monotonic: first run shortly after boot, then every 2 min *after the previous
# deploy finishes* (avoids overlapping docker builds). Wall-clock OnCalendar was
# unreliable on some hosts (stale next_elapse, NEXT shown as "-" in list-timers).
Unit=lunarleague-auto-deploy.service
OnBootSec=45s
OnUnitInactiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
